# docker/gpt2ed/config/gunicorn.conf
[program:gunicorn]
command=/usr/bin/gunicorn app:app -w 9 -b localhost:5001
directory=/deploy/app