# docker/qa/config/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:gunicorn]
# Increase the timeout to 100 seconds.
command=/usr/bin/gunicorn app:app -w 1 -t 100 -b localhost:5001 --log-level info
startsecs=60
directory=/deploy/app
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0